# rainie-station.service
[Unit]
Description=Rainie Station Web Service
After=network.target

[Service]
User=root
# 请修改为服务器上的实际路径
WorkingDirectory=/root/rainie-station
Environment="PATH=/root/rainie-station/venv/bin"
ExecStart=/root/rainie-station/venv/bin/python /root/rainie-station/app.py
Restart=always
RestartSec=2
# 确保 logs 目录存在
StandardOutput=append:/root/rainie-station/logs/web.log
StandardError=append:/root/rainie-station/logs/web.log

[Install]
WantedBy=multi-user.target